# deploy/launchd/com.dahuang.knowledge-auto-update.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dahuang.knowledge-auto-update.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/dahuang/CascadeProjects/knowledge-auto-update/scripts/run_daily_digest.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/dahuang/CascadeProjects/knowledge-auto-update</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/dahuang</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/dahuang/Library/Logs/knowledge-auto-update/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dahuang/Library/Logs/knowledge-auto-update/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
